# server/device-monitor.service
[Unit]
Description=Device Monitor Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/device-monitor/server
ExecStart=/bin/bash run_server.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
